# dummy.service
[Unit]
Description=Dummy Background Service
After=network.target

[Service]
ExecStart=/usr/local/bin/dummy.sh
Restart=always
User=root
StandardOutput=append:/var/log/dummy-service.log
StandardError=append:/var/log/dummy-service.log

[Install]
WantedBy=multi-user.target













